# debian_specs/nginx/nginx-common.postinst.erb
<%
require File.absolute_path(File.dirname(__FILE__)) + '/helpers'
-%>
#!/bin/sh
set -e

# Handle naxsi removal
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi.rules         1.6.2-2~ -- "$@"
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi_core.rules    1.6.2-2~ -- "$@"
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi-ui.conf.1.4.1 1.6.2-2~ -- "$@"
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi-ui.conf       1.6.2-2~ -- "$@"

case "$1" in
  configure)
    logdir="/var/log/nginx"
    # Ensure secure permissions (CVE-2013-0337)
    # http://bugs.debian.org/701112
    #
    # nginx uses 0755 for log files making them world readable,
    # we fix that by using 0750 for the log directory.
    #
    # Allow local admin to override:
    # e.g. dpkg-statoverride --add root adm 0755 /var/log/nginx
    if ! dpkg-statoverride --list "$logdir" >/dev/null; then
      mkdir -p "$logdir"
      chown www-data:adm $logdir
      chmod 0750 $logdir
    fi
<% if is_distribution?("<= saucy") || is_distribution?("<= wheezy") -%>
    # For systemd, enable the nginx unit, but don't start it (just like the
    # init script).
    if [ -x /bin/systemctl ]; then
      systemctl enable nginx.service >/dev/null 2>&1 || true
    fi
<% end -%>
    # If a symlink doesn't exist and can be created, then create it.
    if [ -z $2 ] && [ ! -e /etc/nginx/sites-enabled/default ] &&
       [ -d /etc/nginx/sites-enabled ] && [ -d /etc/nginx/sites-available ]; then
      ln -s /etc/nginx/sites-available/default /etc/nginx/sites-enabled/default
    fi

<% if default_document_root != "/usr/share/nginx/html" -%>
    # Create a default index page when not already present.
    if [ ! -e <%= default_document_root %>/index.nginx-debian.html ]; then
      cp /usr/share/nginx/html/index.html <%= default_document_root %>/index.nginx-debian.html
    fi
<% end -%>

    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

#DEBHELPER#

exit 0
