# debian_specs/nginx/nginx-common.postinst.erb
<%
require File.absolute_path(File.dirname(__FILE__)) + '/helpers'
-%>
#!/bin/sh
set -e

. /usr/share/debconf/confmodule

# Handle naxsi removal
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi.rules         1.6.2-2~ -- "$@"
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi_core.rules    1.6.2-2~ -- "$@"
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi-ui.conf.1.4.1 1.6.2-2~ -- "$@"
dpkg-maintscript-helper rm_conffile \
                 /etc/nginx/naxsi-ui.conf       1.6.2-2~ -- "$@"

case "$1" in
  configure)
    logdir="/var/log/nginx"

    # Allow local admin to override
    if ! dpkg-statoverride --list "$logdir" >/dev/null; then
      mkdir -p "$logdir"
      chown root:adm $logdir
      chmod 0755 $logdir
    fi

    # Secure default logfiles on fresh installations
    if [ -z "$2" ]; then
      access_log="$logdir/access.log"
      error_log="$logdir/error.log"

      if [ ! -e "$access_log" ]; then
        touch "$access_log"
        chmod 640 "$access_log"
        chown www-data:adm "$access_log"
      fi

      if [ ! -e "$error_log" ]; then
        touch "$error_log"
        chmod 640 "$error_log"
        chown www-data:adm "$error_log"
      fi
    fi

<% if is_distribution?("<= saucy") || is_distribution?("<= wheezy") -%>
    # For systemd, enable the nginx unit, but don't start it (just like the
    # init script).
    if [ -x /bin/systemctl ]; then
      systemctl enable nginx.service >/dev/null 2>&1 || true
    fi
<% end -%>
    # If a symlink doesn't exist and can be created, then create it.
    if [ -z $2 ] && [ ! -e /etc/nginx/sites-enabled/default ] &&
       [ -d /etc/nginx/sites-enabled ] && [ -d /etc/nginx/sites-available ]; then
      ln -s /etc/nginx/sites-available/default /etc/nginx/sites-enabled/default
    fi

<% if default_document_root != "/usr/share/nginx/html" -%>
<%
# On Ubuntu 13.04 - 14.10, the default document root is /usr/share/nginx/html.
# This directory also contains an index.html. On newer versions, /usr/share/nginx/html
# still exists but is no longer the default document root; that changed to /var/www/html.
# /var/www/html is empty so we need to create an index page when not already present.
-%>
    # Create a default index page when not already present.
    if [ ! -e <%= default_document_root %>/index.nginx-debian.html ]; then
      cp /usr/share/nginx/html/index.html <%= default_document_root %>/index.nginx-debian.html
    fi
<% end -%>

    ;;

  abort-upgrade|abort-remove|abort-deconfigure)
    ;;

  *)
    echo "postinst called with unknown argument \`$1'" >&2
    exit 1
    ;;
esac

#DEBHELPER#

exit 0
